# runit/postfix/run
#!/bin/sh

postalias /etc/postfix/aliases

exec postfix start-fg
